# templates/odoo.service.erb
[Unit]
Description=Advanced OpenSource ERP and CRM server
Requires=postgresql.service
After=postgresql.service

[Service]
User="<%= scope['odoo::params::odoo_user'] %>"
Group="<%= scope['odoo::params::odoo_group'] %>"
Type=simple
PermissionsStartOnly=true
EnvironmentFile=-/etc/conf.d/openerp-server
SyslogIdentifier=odoo.tietoteema.com
PIDFile=/run/odoo/server.pid
ExecStartPre=/usr/bin/install -d -m755 -o "<%= scope['odoo::params::odoo_user'] %>" -g "<%= scope['odoo::params::odoo_group'] %>" /run/openerp
ExecStart=/opt/odoo/odoo-bin -c /etc/odoo.conf --unaccent --without-demo=True --no-database-list
ExecStop=/bin/kill $MAINPID

[Install]
WantedBy=multi-user.target
