# templates/odoo.service.erb
[Unit]
Description=Advanced OpenSource ERP and CRM server

[Service]
User=<%= scope['odoo::odoo_user'] %>
Group=<%= scope['odoo::odoo_group'] %>
Type=simple
SyslogIdentifier=odoo
PIDFile=<%= scope['odoo::pidpath'] %>/server.pid
ExecStart=<%= scope['odoo::odoo_executable'] %> -c /etc/odoo.conf --no-database-list --db-filter=<%= scope['odoo::db_filter'] %>
ExecStop=/bin/kill $MAINPID

[Install]
WantedBy=multi-user.target
